# email-daemon.service
[Unit]
Description=Email AI Daemon
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ws/MailMind

# 环境变量（按需修改 /etc/mailmind.env）
EnvironmentFile=/etc/mailmind.env

# 启动命令（按需修改 --mailbox 和 --ai）
ExecStart=/home/ubuntu/ws/MailMind/venv/bin/python email_daemon.py --mailbox 126 --ai claude

# 崩溃自动重启
Restart=on-failure
RestartSec=10

# 日志
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
