# email-daemon.service
[Unit]
Description=Email AI Daemon
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ws/ccmail

# 环境变量（按需修改）
Environment=MAIL_126_ADDRESS=your@126.com
Environment=MAIL_126_PASSWORD=your-auth-code
Environment=MAIL_126_ALLOWED=your@126.com
Environment=POLL_INTERVAL=60

# 启动命令（按需修改 --mailbox 和 --ai）
ExecStart=/home/ubuntu/ws/ccmail/venv/bin/python email_daemon.py --mailbox 126 --ai claude

# 崩溃自动重启
Restart=on-failure
RestartSec=10

# 日志
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
